# deployment/systemd/asia-frontend.service
[Unit]
Description=Asia Site Next.js Frontend
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/asia-site/frontend
Environment="NODE_ENV=production"
Environment="PORT=3010"
EnvironmentFile=/var/www/asia-site/frontend/.env.local
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=asia-frontend

[Install]
WantedBy=multi-user.target
