# deployment/systemd/asia-backend.service
[Unit]
Description=Asia Site Backend (Django/Gunicorn)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/asia-site/backend
Environment="PATH=/var/www/asia-site/backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/var/www/asia-site/backend/venv/bin/gunicorn \
    --workers 3 \
    --bind 127.0.0.1:8000 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile - \
    --log-level info \
    config.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=asia-backend

[Install]
WantedBy=multi-user.target
